# common/overlays/rootfs/frecon/etc/profile.d/frecon.sh
export FRECON_SHELL=bash
export FRECON_VTS=1
# export FRECON_VT1=1
export FRECON_FB_ROTATE=0
export FRECON_FB_SCALE=1
